# docker.conf
[program:docker]
command=docker run -v /opt/run/webphp/config:/etc -v /opt/run/webphp/log:/var/log/apache2 zadatak

