# docker.conf
[program:docker]
command=docker run -v /opt/run/webphp/config/apache2:/etc/apache2 -v /opt/run/webphp/config/php:/etc/php -v /opt/run/webphp/log:/var/log/apache2 -p 47280:80 zadatak
autostart=true

